# tmp/nuts-cc-file.list
/etc/alternatives/README
/etc/apt/
/etc/calamares/
/etc/cron.daily/dpkg
/etc/dpkg/
/etc/group-
/etc/gshadow-
/etc/logrotate.d/alternatives
/etc/logrotate.d/dpkg
/etc/shadow-
/etc/skel/Desktop/org.nxos.installer.desktop
/etc/subgid-
/etc/subuid-
/etc/xdg/autostart/calamares-kwin-rules.desktop
/opt/1
/opt/2
/opt/3
/usr/bin/calamares
/usr/bin/dpkg
/usr/bin/dpkg-*
/usr/bin/update-alternatives
/usr/include/libcalamares/
/usr/lib/x86_64-linux-gnu/calamares/
/usr/libexec/dpkg
/usr/share/applications/calamares.desktop
/usr/share/calamares/
/usr/share/doc/calamares/
/usr/share/doc/dpkg/
/usr/share/dpkg/
/usr/share/icons/hicolor/scalable/apps/calamares.svg
/usr/share/lintian/overrides/calamares
/usr/share/lintian/overrides/dpkg
/usr/share/man/man1/calamares.1.gz
/usr/share/man/man8/calamares.8.gz
/usr/share/polkit-1/actions/com.github.calamares.calamares.policy
/usr/share/polkit-1/actions/org.dpkg.pkexec.update-alternatives.policy
/var/cache/apt
/var/lib/apt/lists
/var/lib/dpkg/
